# dot_config/autostart/corectrl.desktop
[Desktop Entry]
Exec=corectrl
Name=corectrl
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
