# dot_config/autostart/Spike.desktop
[Desktop Entry]
Exec=firefoxpwa site launch 01FG7YSGVCFA4CM015VZ0F527A
Name=Spike
Type=Application
Version=1.0
