# dot_config/autostart/tint2.desktop
[Desktop Entry]
Exec=tint2
Name=tint2
Type=Application
Version=1.0
