# dot_config/autostart/xmodmap.desktop
[Desktop Entry]
Exec=xmodmap ~/.Xmodmap
Name=xmodmap
Type=Application
Version=1.0
